# vslam_libs/sba/build/devel/etc/catkin/profile.d/05.catkin-test-results.sh
# generated from catkin/cmake/env-hooks/05.catkin-test-results.sh.develspace.em

export CATKIN_TEST_RESULTS_DIR="/home/sebasgm/ros_devel/src/slam_karto/vslam_libs/sba/build/test_results"
export ROS_TEST_RESULTS_DIR=$CATKIN_TEST_RESULTS_DIR
